# .config/systemd/user/teams-jiggler-off.service
[Unit]
Description=Switch Teams jiggler off

[Service]
Type=oneshot
ExecStart=%h/.local/bin/teams-jiggler-toggle off
