# .config/systemd/user/teams-jiggler-schedule.service
[Unit]
Description=Switch Teams jiggler to scheduled mode

[Service]
Type=oneshot
ExecStart=%h/.local/bin/teams-jiggler-toggle schedule
