# .config/autostart/bitwarden.desktop
[Desktop Entry]
  Type=Application
  Version=2026.6.1
  Name=Bitwarden
  Comment=Bitwarden startup script
  Exec=/opt/Bitwarden/bitwarden-app --autostart
  StartupNotify=false
  Terminal=false